# server.ini
[program:server]
directory=/usr/local/services/server/
command=npm run start
autostart=true
autorestart=true
startsecs=10
startretries=100
